# instance/systemd/docker-git.service
[Unit]
Description=Git container
Requires=setup-instance.service docker.service
After=setup-instance.service docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker start --attach git
ExecStop=/usr/bin/docker stop git