# instance/systemd/group-git.service
[Unit]
Description=Create Git group
ConditionFirstBoot=true

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/groupadd --gid 1001 git