# instance/systemd/media-volume.mount
[Unit]
Description=Mount volume to /media/volume
Requires=dev-vdb.device
After=dev-vdb.device

[Mount]
What=/dev/vdb1
Where=/media/volume
Type=ext4